# rag/openclaw-rag.service
[Unit]
Description=OpenClaw RAG Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/alireza/.openclaw/workspace/rag
Environment=RAG_DATA_DIR=/home/alireza/.openclaw/rag
Environment=RAG_EMBED_MODEL=BAAI/bge-small-en-v1.5
ExecStart=/home/alireza/.openclaw/rag/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8799
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
